# my_docker_config/root_overlay/etc/systemd/system/mnt-ssd2.mount
[Unit]
Description=mounting Storage2 to SSD2

[Mount]
Where=/mnt/ssd2
What=/dev/disk/by-label/storage2

[Install]
WantedBy=multi-user.target
