# my_docker_config/root_overlay/etc/systemd/system/mnt-ssd1.mount
[Unit]
Description=mounting Storage1 to SSD1

[Mount]
Where=/mnt/ssd1
What=/dev/disk/by-label/storage1

[Install]
WantedBy=multi-user.target
